# service/webp.service
[Unit]
Description=Webp Converter

[Service]
ExecStart=/home/vxphong/scripts/monitor-file-creation.sh

[Install]
WantedBy=multi-user.target
